# .config/systemd/user/clash.service
[Unit]
Description=fuckGFW

[Service]
Type=simple
ExecStart=/home/jason/.nix-profile/bin/clash

[Install]
WantedBy=default.target

