# .config/systemd/user/clash.service
[Unit]
Description=fuckGFW

[Service]
Type=simple
ExecStart=/usr/bin/clash

[Install]
WantedBy=multi-user.target

